# roles/envoy/templates/envoy.service

[Service]
ExecStart=/usr/local/bin/envoy -c /etc/envoy/envoy.yaml 
Restart=always
StandardOutput=journal

[Install]
WantedBy=multi-user.target